# finagent.service
[Unit]
Description=FinAgent Backend API
After=network.target

[Service]
Type=simple
User=mk
WorkingDirectory=/app/workspace/finAgent
Environment="PATH=/app/workspace/finAgent/finagent/bin:/usr/local/bin:/usr/bin:/bin"

# Environment variables
EnvironmentFile=-/etc/finagent/env
EnvironmentFile=-/etc/environment

ExecStart=/app/workspace/finAgent/finagent/bin/gunicorn -w 2 -k uvicorn.workers.UvicornWorker finagent_api:app --bind 0.0.0.0:8000 --timeout 480 --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
