# finagent.service
[Unit]
Description=FinAgent Backend API
After=network.target

[Service]
Type=simple
User=mk
WorkingDirectory=/app/workspace/finAgent
Environment="PATH=/app/workspace/finAgent/finagent/bin:/usr/local/bin:/usr/bin:/bin"

# Load environment variables from file (API keys, etc.)
# Create this file with: sudo ./setup-service-env.sh
EnvironmentFile=/etc/finagent/env

ExecStart=/app/workspace/finAgent/finagent/bin/gunicorn -w 2 -k uvicorn.workers.UvicornWorker finagent_api:app --bind 0.0.0.0:8000 --timeout 480 --log-level info
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
